# wrapper/systemd/ssh-rt-inner-sshd.service
[Unit]
Description=ssh-rt-auth inner sshd (hermetic, localhost-only)
BindsTo=ssh-rt-wrapperd.service
After=ssh-rt-wrapperd.service

[Service]
Type=forking
ExecStart=/usr/sbin/sshd -f /var/lib/ssh-rt-auth/inner-sshd/sshd_config
PIDFile=/var/lib/ssh-rt-auth/inner-sshd/sshd.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
