# systemd/msshd.service
[Unit]
Description=msshd — ssh-rt-auth wrapper (mTLS-authenticated SSH frontend)
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=_msshwrapper
Group=_msshwrapper
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/local/bin/msshd --config /etc/ssh-rt-auth/wrapper.yaml
Restart=on-failure
RestartSec=5

# Sandboxing
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/var/log/ssh-rt-auth /var/lib/ssh-rt-auth
ProtectHome=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
